# kmonad/kmonad@.service
[Unit]
BindTo=%i.device
After=%i.device

[Service]
Type=simple
ExecStart=/usr/bin/kmonad-autoconfig "/%I"
